# backend/com.tinyapp.DBngin.postgresql-40AF4CDE-8BFD-4EB6-B1CD-C37A9DD5387F.plist
  <?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>com.tinyapp.DBngin.postgresql-40AF4CDE-8BFD-4EB6-B1CD-C37A9DD5387F</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/Shared/DBngin/postgresql/17.0/bin/postgres</string>
      <string>-D</string>
      <string>/Users/admin/Projects/RatelShop/backend</string>
      <string>-p</string>
      <string>5432</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/Shared/DBngin/postgresql/17.0</string>
    <key>StandardOutPath</key>
    <string>/Users/admin/Projects/RatelShop/backend/postgres_out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/Projects/RatelShop/backend/postgres_err.log</string>
  </dict>
  </plist>